# deploy/systemd/datavis.service
[Unit]
Description=datavis.au web app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/cTrade
EnvironmentFile=/home/ec2-user/cTrade/.env
ExecStart=/home/ec2-user/cTrade/venv/bin/uvicorn datavis.app:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
