# deploy/systemd/engineering-orchestrator.service
[Unit]
Description=datavis bounded engineering orchestrator
After=network-online.target research-control.service
Wants=network-online.target research-control.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/cTrade
EnvironmentFile=-/etc/datavis.env
EnvironmentFile=-/etc/datavis-research.env
EnvironmentFile=-/etc/datavis-control.env
ExecStart=/home/ec2-user/venvs/datavis/bin/python -u -m datavis.control.engineering_orchestrator_runtime
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

