# deploy/systemd/ottprocessor.service
[Unit]
Description=datavis OTT incremental processor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/cTrade
EnvironmentFile=/etc/datavis.env
ExecStart=/home/ec2-user/venvs/datavis/bin/python -m datavis.ott_jobs worker --symbol XAUUSD --source mid --matype VAR --length 2 --percent 1.4 --poll-seconds 1.0 --batch-size 2000 --log-every 20
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
