# deploy/systemd/unity.service
[Unit]
Description=cTrade UNITY Processor (Candidate Birth + Paper Journal)
After=network-online.target postgresql.service tickcalc.service
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/cTrade
ExecStart=/home/ec2-user/cTrade/venv/bin/python -u /home/ec2-user/cTrade/unityFromDB.py --mode live --symbol XAUUSD
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
